# Dryft/infra/dryft-api.service
[Unit]
# Human-readable description shown by systemctl status.
Description=Dryft API service
# Ensure networking is ready before startup.
After=network-online.target
Wants=network-online.target

[Service]
# Run the service binary directly.
Type=simple
# Run with least-privilege dedicated service account.
User=dryft
Group=dryft
# All relative paths resolve from the deployment directory.
WorkingDirectory=/opt/dryft
# Load environment variables from production env file.
EnvironmentFile=/opt/dryft/.env.prod
# Primary process entrypoint.
ExecStart=/opt/dryft/dryft-api
# Restart on crashes/non-zero exit codes.
Restart=on-failure
# Wait briefly before restart attempts.
RestartSec=5
# Send logs to journald.
StandardOutput=journal
StandardError=journal
# Service hardening: mount system paths read-only for this unit.
ProtectSystem=strict
# Block access to other users' home directories.
ProtectHome=true
# Prevent privilege escalation (e.g., setuid binaries).
NoNewPrivileges=true
# Explicit writable path needed by the service.
ReadWritePaths=/opt/dryft

[Install]
# Start on normal multi-user boot.
WantedBy=multi-user.target

